# roles/roq-gemini/templates/roq-gemini.timer
[Unit]
Description=gemini gateway timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
