# roles/samples/templates/roq-collector.timer
[Unit]
Description=Market data collector start timer

[Timer]
OnCalendar=Mon-Fri *-*-* 09:14:00

[Install]
WantedBy=multi-user.target
