# templates/gateway/item.timer
[Unit]
Description=roq-{{ config.options.name }} timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
