# roles/roq-bitmex/templates/roq-bitmex-backup.timer
[Unit]
Description=bitmex gateway backup timer

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true

[Install]
WantedBy=multi-user.target
