# roles/roq-bitmex/templates/roq-bitmex.timer
[Unit]
Description=bitmex gateway timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
