# roles/roq-benchmark/templates/roq-ping.service
[Unit]
Description=Ping {{ item }}
After=network-online.target

[Service]
UMask=0000
User={{ roq_user if roq_user is defined else ansible_user_id }}
Group={{ roq_user if roq_user is defined else ansible_user_id }}
ExecStart={{ root }}/usr/local/bin/roq-ping-{{ item }}.sh
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
