# roles/grafana/templates/grafana-backup.timer
[Unit]
Description=Grafana backup timer

[Timer]
OnCalendar=*-*-* 22:45:00
Persistent=true

[Install]
WantedBy=multi-user.target
