# roles/roq-coinbase-pro/templates/roq-coinbase-pro-stop.timer
[Unit]
Description=coinbase pro gateway stop timer

[Timer]
OnCalendar={{ config.service.stop }}

[Install]
WantedBy=multi-user.target
