# roles/roq-coinbase-pro/templates/roq-coinbase-pro-stop.timer
[Unit]
Description=Coinbase-Pro stop timer

[Timer]
OnCalendar=Mon-Sun *-*-* 23:59:00

[Install]
WantedBy=multi-user.target
