# roles/roq-coinbase-pro/templates/roq-coinbase-pro.timer
[Unit]
Description=coinbase pro gateway timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
