# roles/femas/templates/roq-femas-stop.timer
[Unit]
Description=Femas gateway stop timer

[Timer]
OnCalendar=Mon-Fri *-*-* 15:16:00

[Install]
WantedBy=multi-user.target
