# roles/roq-kraken/templates/roq-kraken.timer
[Unit]
Description=kraken gateway timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
