# roles/roq-deribit/templates/roq-deribit-stop.timer
[Unit]
Description=Deribit stop timer

[Timer]
OnCalendar=Mon-Sun *-*-* 23:59:00

[Install]
WantedBy=multi-user.target
