# roles/roq-deribit/templates/roq-deribit.timer
[Unit]
Description=Deribit start timer

[Timer]
OnCalendar=Mon-Sun *-*-* 00:01:00

[Install]
WantedBy=multi-user.target
