# roles/certbot/templates/certbot.timer
[Unit]
Description=certbot renew timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
